# powertop.service
# PowerTOP auto-tune for battery optimization.
# Applies power-saving tunings on boot, particularly useful
# for MacBook hardware running on battery.
# Ref: https://umarahmad.xyz/blog/installing-fedora-on-macbook/

[Unit]
Description=PowerTOP auto-tune for battery optimization

[Service]
Type=oneshot
ExecStart=/usr/sbin/powertop --auto-tune

[Install]
WantedBy=multi-user.target
